# .crontab
# install me using "crontab .crontab"!
#
*/20 * * * * check-battery.sh
